# files/client/check-certificate-existence@.timer
[Unit]
Description=Check x.509 certificate existence

[Timer]
OnActiveSec=1
OnCalendar=daily

[Install]
WantedBy=multi-user.target
